# agent/systemd/happydo-agent.service
[Unit]
Description=HappyDo Guard RTSP-to-RTMP Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/happydo/rtsp-to-rtmp.sh /etc/happydo-agent.conf
Restart=always
RestartSec=15
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
